# winspotlightkde.service
[Unit]
Description=WinSpotlightKDE
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/winspotlightkde
