# winspotlightkde.service
[Unit]
Description=WinSpotlightKDE

[Service]
Type=simple
ExecStart=/opt/winspotlightkde